# root/etc/services.d/rtorrent/run
#!/usr/bin/with-contenv sh

exec supervisord -c /etc/supervisord.conf -n
